# .devcontainer/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=info

[program:start-desktop]
command=/usr/local/bin/start-desktop.sh
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/start-desktop.log
redirect_stderr=true
